# speedtest-mqtt/rootfs/etc/s6-overlay/s6-rc.d/speedtest/run
#!/usr/bin/with-contenv bashio
# s6-overlay service entry point
exec /usr/bin/speedtest_loop.sh
